# .config/systemd/user/album-art.timer
[Unit]
Description=Run album art updater every second

[Timer]
OnBootSec=1
OnUnitActiveSec=1
AccuracySec=200ms

[Install]
WantedBy=default.target
