# .config/systemd/user/matugen.service
[Unit]
Description=Run Matugen to generate a new wallpaper

[Service]
Type=oneshot
ExecStart=/usr/bin/update-wallpaper
