# .config/systemd/user/wellbeing-monitor.service
[Unit]
Description=Digital Wellbeing Battery Monitor
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/wellbeing monitor
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
