# .config/systemd/user/xdg-desktop-portal-hyprland.service.d/override.conf
[Service]
Environment="QT_PLUGIN_PATH=/usr/lib/qt6/plugins"
Environment="QT_QPA_PLATFORM=wayland;xcb"
